# supervisord/reusable/clamav.conf
[program:clamav-update]
command=/usr/bin/freshclam -d -l %(ENV_HOME)s/logs/freshclam.log --foreground=true
numprocs=1
startsecs=0
