# supervisord/reusable/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
numprocs=1
process_name=%(program_name)s-%(process_num)d
