# supervisord/clamav_reusable.conf
[program:clamav-update]
command=bash -c "if [ \"${CLAM_AV}\" = 'yes' ]; then /usr/bin/freshclam -d \
    -l %(ENV_HOME)s/logs/freshclam.log --foreground=true; fi"
numprocs=1
startsecs=0
